# z-devel/perforce-server/crontab
#
# Verify and backup the Perforce repository daily
#
05 03 * * * /home/perforce/scripts/backup.sh

