# z-mail/pine-pop2/xinetd.d/pop2
service pop2
{
	socket_type		= stream
	protocol		= tcp
	wait			= no
	user			= root
	server			= /usr/sbin/ipop2d
#	server_args		= 
	per_source		= 4
	log_type		= FILE /var/log/xinetd/pop2
	log_on_success		= PID HOST USERID EXIT DURATION
	log_on_failure		= USERID ATTEMPT
}
